# service/openrc
#!/sbin/openrc-run

# Config
USER='discord'
DIR='anchorAnalytics'

depend() {
    need net
}

procname='discord-bot-anchor-analytics'
pidfile="/run/$procname.pid"

command_user="$USER:$USER"
command_background='true'

command='python3'
command_args="/home/$USER/$DIR/bot.py"
