# service/openrc
#!/sbin/openrc-run

USER=discord

DIR=anchorAnalytics
DIR=/home/$USER/$DIR

EXEC=bot.py
EXEC=$DIR/$EXEC

ID=discord-bot-anchor-analytics
ID=/run/$ID.pid

command='start-stop-daemon'
command_args="\
--start \
--user $USER \
--group $USER \
--chdir $DIR \
--make-pidfile \
--pidfile $ID \
--background \
--stdout $DIR/out.txt \
--stderr $DIR/err.txt \
--exec $EXEC"

depend() {
    need net
}
