# service/openrc
#!/sbin/openrc-run

USER='discord'

command='/usr/bin/env'
command_args="PYTHONPATH=/home/$USER/anchorAnalytics/dep /usr/bin/python3 /home/$USER/anchorAnalytics/bot.py"
command_user="$USER:$USER"
description='Discord Anchor'
supervisor='supervise-daemon'

depend() {
	need net
	after firewall
	use dns
}
