# service/openrc
#!/sbin/openrc-run

USER=discord

DIR=anchorAnalytics
DIR=/home/$USER/$DIR

EXEC=bot.py
EXEC=$DIR/$EXEC

ID=discord-bot-anchor-analytics
ID=/run/$ID.pid

SUPERVISOR=start-stop-daemon

START=" \
	--start \
	--user $USER \
	--group $USER \
	--chdir $DIR \
	--make-pidfile \
	--pidfile $ID \
	--background \
	--stdout $DIR/out.txt \
	--stderr $DIR/err.txt \
	--exec $EXEC"

STOP=" \
	--stop \
	--pidfile $ID"

description='Discord Bot - Anchor Analytics'

depend() {
	need net
	after firewall
	use dns
}

start() {
	$SUPERVISOR $START
}

stop() {
	$SUPERVISOR $STOP
}
